# scribe/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:scribe]
command=/usr/local/bin/scribed -c /usr/local/etc/scribe_example.conf
